# platform/centec-arm64/sonic-platform-modules-e550/24t16y/service/24t16y_platform.service
[Unit]
Description=Centec modules init
After=local-fs.target
Before=syncd.service

[Service]
Type=oneshot
ExecStart=-/etc/init.d/platform-modules-e550-24t16y start
ExecStop=-/etc/init.d/platform-modules-e550-24t16y stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
