# platform/centec-arm64/sonic-platform-modules-e550/debian/platform-modules-e550-24t16y.postinst
systemctl enable 24t16y_platform.service
systemctl start 24t16y_platform.service
