# platform/centec-arm64/sonic-platform-modules-e550/debian/platform-modules-e550-24x8y2c.postinst
systemctl enable 24x8y2c_platform.service
systemctl start 24x8y2c_platform.service
