# platform/centec-arm64/platform-conf/debian/centec-platform-conf.postinst
systemctl enable centec_platform_conf.service
systemctl start centec_platform_conf.service
